# services/stdatalog-cli.service
[Unit]
Description=STDatalog CLI Logger Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=kirwinr
Group=kirwinr
WorkingDirectory=/home/kirwinr/Desktop/stdatalog/stdatalog_examples/gui_applications/stdatalog/CLI
Environment=PATH=/home/kirwinr/Desktop/stdatalog/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=VIRTUAL_ENV=/home/kirwinr/Desktop/stdatalog/.venv
ExecStart=/home/kirwinr/Desktop/stdatalog/.venv/bin/python stdatalog_CLI.py
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3
StandardOutput=append:/home/kirwinr/logs/stdatalog-cli.log
StandardError=append:/home/kirwinr/logs/stdatalog-cli.log
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
