# services/stdatalog-usboffload.service
[Unit]
Description=STDatalog USB Acquisition Data Offload Service
After=network.target

[Service]
Type=simple
# Update this path to your Python interpreter in your venv!
ExecStart=/home/kirwinr/Desktop/stdatalog/.venv/bin/python /home/kirwinr/Desktop/stdatalog/thread/usb_transfer.py
Environment=PYTHONPATH=/home/kirwinr/Desktop/stdatalog:/home/kirwinr/Desktop/stdatalog/.venv/lib/python3.11/site-packages
WorkingDirectory=/home/kirwinr/Desktop/stdatalog/thread
User=kirwinr
Group=kirwinr
Restart=always
RestartSec=5
StandardOutput=append:/home/kirwinr/logs/stdatalog-usboffload.log
StandardError=append:/home/kirwinr/logs/stdatalog-usboffload.log

[Install]
# No WantedBy target - this service must be started manually
# It will restart on failure but won't auto-start on boot
